# upsplus.timer
[Unit]
Description=run upsplus

[Timer]
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
